# docker/supervisord.conf
# ─── Supervisord config for Clever Cloud production (Ubuntu 24.04) ─────────────
# Manages three processes: nginx (proxy), backend (FastAPI), frontend (Next.js)
# All logs go to stdout/stderr so Clever Cloud captures them.
#
# NOTE: No 'environment=' directives — supervisord inherits the full container
# environment by default, which already contains all Clever Cloud injected vars.

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

# ── nginx: reverse proxy (port 8080) ──────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── FastAPI backend (port 8000, internal only) ─────────────────────────────────
[program:backend]
directory=/app/backend
command=/app/.venv/bin/gunicorn app.main:app --workers 4 --worker-class uvicorn.workers.UvicornWorker --bind 127.0.0.1:8000 --timeout 120 --access-logfile - --error-logfile - --log-level info
autostart=true
autorestart=true
priority=20
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Next.js frontend (port 3000, internal only) ────────────────────────────────
[program:frontend]
directory=/app/frontend-standalone
command=node server.js
autostart=true
autorestart=true
priority=20
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
