# docker/supervisord.conf
# ─── Supervisord config for Clever Cloud production (Ubuntu 24.04) ─────────────
# Manages three processes: nginx (proxy), backend (FastAPI), frontend (Next.js)
# All logs go to stdout/stderr so Clever Cloud captures them.

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

# ── nginx: reverse proxy (port 8080) ──────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── FastAPI backend (port 8000, internal only) ─────────────────────────────────
[program:backend]
directory=/app/backend
command=uv run gunicorn app.main:app
    --workers 4
    --worker-class uvicorn.workers.UvicornWorker
    --bind 127.0.0.1:8000
    --timeout 120
    --access-logfile -
    --error-logfile -
    --log-level info
autostart=true
autorestart=true
priority=20
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Pass all environment variables through to the backend process
environment=
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    REDIS_URL="%(ENV_REDIS_URL)s",
    CG_ENV="%(ENV_CG_ENV)s",
    CG_HTTP_PORT="%(ENV_CG_HTTP_PORT)s",
    CG_LOG_LEVEL="%(ENV_CG_LOG_LEVEL)s",
    CG_MASTER_KEY="%(ENV_CG_MASTER_KEY)s",
    CG_JWT_SECRET="%(ENV_CG_JWT_SECRET)s",
    CG_JWT_ALGORITHM="%(ENV_CG_JWT_ALGORITHM)s",
    CG_JWT_ACCESS_EXPIRE_MINUTES="%(ENV_CG_JWT_ACCESS_EXPIRE_MINUTES)s",
    CG_JWT_REFRESH_EXPIRE_DAYS="%(ENV_CG_JWT_REFRESH_EXPIRE_DAYS)s",
    CG_ENCRYPTION_KEY="%(ENV_CG_ENCRYPTION_KEY)s",
    CG_ADMIN_USERNAME="%(ENV_CG_ADMIN_USERNAME)s",
    CG_ADMIN_PASSWORD="%(ENV_CG_ADMIN_PASSWORD)s",
    CG_ADMIN_EMAIL="%(ENV_CG_ADMIN_EMAIL)s",
    CG_ADMIN_FIRST_NAME="%(ENV_CG_ADMIN_FIRST_NAME)s",
    CG_ADMIN_LAST_NAME="%(ENV_CG_ADMIN_LAST_NAME)s",
    CG_ROUTING_STRATEGY="%(ENV_CG_ROUTING_STRATEGY)s",
    CG_NUM_RETRIES="%(ENV_CG_NUM_RETRIES)s",
    CG_TIMEOUT="%(ENV_CG_TIMEOUT)s",
    CG_ALLOWED_FAILS="%(ENV_CG_ALLOWED_FAILS)s",
    CG_COOLDOWN_TIME="%(ENV_CG_COOLDOWN_TIME)s",
    CG_CORS_ORIGINS="%(ENV_CG_CORS_ORIGINS)s"

# ── Next.js frontend (port 3000, internal only) ────────────────────────────────
[program:frontend]
directory=/app/frontend-standalone
command=node server.js
autostart=true
autorestart=true
priority=20
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    PORT="3000",
    HOSTNAME="127.0.0.1",
    NODE_ENV="production",
    NEXT_PUBLIC_API_URL="%(ENV_NEXT_PUBLIC_API_URL)s"
